# backend/backend/supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=uvicorn server:app --host 0.0.0.0 --port 8000

[program:dramatiq]
command=dramatiq background_tasks
